# scripts/workout-deploy.service.template
[Unit]
Description=Workout App Deploy Script
After=network.target
Wants=network.target

[Service]
Type=oneshot
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{WORKING_DIR}}
Environment=PATH={{PATH}}
Environment=HOME={{HOME}}
ExecStart={{WORKING_DIR}}/scripts/deploy.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
