# scripts/workout-deploy.service
[Unit]
Description=Workout App Deploy Script
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=pepper
Group=users
WorkingDirectory=/var/services/homes/pepper/workout
Environment=PATH=/var/services/homes/pepper/.nvm/versions/node/v22.11.0/bin:/sbin:/bin:/usr/sbin:/usr/bin:/usr/syno/sbin:/usr/syno/bin:/usr/local/sbin:/usr/local/bin
Environment=HOME=/var/services/homes/pepper
ExecStart=/var/services/homes/pepper/workout/scripts/deploy.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
